# sensu-server/files/supervisord/152-sensu-client.conf
[program:sensu-client]
command=/etc/init.d/sensu-client start
