# sensu-server/files/supervisord/153-uchiwa.conf
[program:uchiwa]
command=/etc/init.d/uchiwa start
